# templates/lamp/containers/0/php/crontab
59  23  *   *   *   root /usr/local/bin/php /var/www/drakeroll.com/bin/console game:hash:close
